# crontab
0 21 * * * /usr/bin/python3 /usr/app/src/dellStartBackup.py
0 2 * * * /usr/bin/python3 /usr/app/src/dellStartBackup.py
